# desktop/DEBIAN/postinst
#!/bin/bash
 
if [ "$1" = configure ]; then
   systemctl enable soundqube-desktop
   systemctl daemon-reload
   systemctl start soundqube-desktop
fi
