# fluidsynth/DEBIAN/postinst
#!/bin/bash
 
if [ "$1" = configure ]; then
   systemctl disable --global fluidsynth.service
fi
