# fluidsynth/DEBIAN/postinst
#!/bin/bash
 
if [ "$1" = configure ]; then
   systemctl enable soundqube-fluidsynth
   systemctl daemon-reload
   systemctl start soundqube-fluidsynth
fi
