# fluidsynth/DEBIAN/prerm
#!/bin/bash
 
if [ "$1" = remove ] ; then
        systemctl stop soundqube-fluidsynth
        systemctl disable soundqube-fluidsynth
fi
