# system/etc/systemd/user/soundqube-desktop@.service
[Unit]
Description=Desktop for screen %i

[Service]
ExecStart=/usr/bin/blackbox -display :%i
Restart=on-failure

[Install]
WantedBy=default.target
