# surge/DEBIAN/postinst
#!/bin/bash
 
if [ "$1" = configure ]; then
   systemctl enable soundqube-surge
   systemctl daemon-reload
   systemctl start soundqube-surge
fi
