# setbfree/DEBIAN/prerm
#!/bin/bash
 
if [ "$1" = remove ] ; then
        systemctl stop soundqube-setbfree
        systemctl disable soundqube-setbfree
fi
