# system/DEBIAN/prerm
#!/bin/bash
 
if [ "$1" = remove ] ; then
        systemctl stop media-soundqube.mount
        systemctl disable media-soundqube.mount
fi

