# autoconnect/DEBIAN/postinst
#!/bin/bash
 
if [ "$1" = configure ]; then
   systemctl enable soundqube-autoconnect
   systemctl daemon-reload
   systemctl stop soundqube-autoconnect
   systemctl start soundqube-autoconnect
fi
