# autoconnect/DEBIAN/prerm
#!/bin/bash
 
if [ "$1" = remove ] ; then
        systemctl stop soundqube-autoconnect
        systemctl disable soundqube-autoconnect
fi
