# fluidweb/DEBIAN/postinst
#!/bin/bash
 
if [ "$1" = configure ]; then
   systemctl enable fluidweb
   systemctl daemon-reload
   systemctl start fluidweb
fi
