# fluidweb/DEBIAN/prerm
#!/bin/bash
 
if [ "$1" = remove ] ; then
        systemctl stop fluidweb
        systemctl disable fluidweb
fi
