# unison/DEBIAN/postinst
#!/bin/bash
 
if [ "$1" = configure ]; then
   systemctl enable soundqube-unison
   systemctl daemon-reload
   systemctl start soundqube-unison
fi
