# unison/DEBIAN/prerm
#!/bin/bash
 
if [ "$1" = remove ] ; then
        systemctl stop soundqube-unison
        systemctl disable soundqube-unison
fi
