# aeolus/DEBIAN/prerm
#!/bin/bash
 
if [ "$1" = remove ] ; then
        systemctl stop soundqube-aeolus
        systemctl disable soundqube-aeoljm
fi
