# jack/DEBIAN/postinst
#!/bin/bash
 
if [ "$1" = configure ]; then
   systemctl enable soundqube-jack
   systemctl daemon-reload
   systemctl start soundqube-jack
fi
